# browser_worker/deploy/novnc.service
[Unit]
Description=noVNC Browser-Based VNC Client
After=x11vnc.service
Requires=x11vnc.service

[Service]
Type=simple
User=root
ExecStart=/usr/share/novnc/utils/novnc_proxy \
    --vnc localhost:5900 \
    --listen 0.0.0.0:6080 \
    --idle-timeout 0
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
